# custom/rootfs/etc/crontab
*/5 * * * * /root/createCloudLogFile.sh
